# features/services/3x-ui/templates/backup-3x-ui.service.j2
[Unit]
Description=Backup 3x-ui data

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup-3x-ui
